# aports/device-motorola-titan/50-firmware.rules
# do not edit this file, it will be overwritten on update
#SUBSYSTEM=="firmware", ACTION=="add", RUN{builtin}="firmware"

# https://github.com/postmarketOS/pmbootstrap/issues/143
# /usr/local not allowed in Alpine apks
SUBSYSTEM=="firmware", ACTION=="add", RUN+="/usr/bin/firmware.sh"
